# .config/systemd/user/lspmux.service
[Unit]
Description=Language server multiplexer (lspmux) — user service
Documentation=https://codeberg.org/p2502/lspmux

[Service]
Type=simple
# Prefer cargo-installed lspmux if present; fall back to common system paths via PATH
Environment=PATH=%h/.cargo/bin:/usr/local/bin:/usr/bin
RuntimeDirectory=lspmux
RuntimeDirectoryMode=0700
ExecStartPre=/usr/bin/rm -f %t/lspmux/lspmux.sock
ExecStart=/usr/bin/env lspmux server
Restart=on-failure
RestartSec=1s

[Install]
WantedBy=default.target
